# MskKanbanBoard.service.in
[Unit]
Description=FormalisedLog web service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=formalisedlog
WorkingDirectory=${APPROOT} 
ExecStart=/usr/bin/env FormalisedLog.wt --docroot ${APPROOT} --approot ${APPROOT} --http-port 8085 --http-addr 127.0.0.1 --accesslog /var/lib/FormalisedLog/access.log

[Install]
WantedBy=multi-user.target
